# OSimages/SlitazOS/var/lib/tazpkg/installed/slitaz-base-files/files.list
/sbin/mktazdevs.sh
/usr/bin/calc
/usr/bin/ldd
/usr/bin/man
/usr/lib/slitaz/httphelper
/usr/share/applications/httpd-cgi.desktop
/var/www/cgi-bin/cgi-env.sh
/var/www/cgi-bin/index.cgi
/etc/adduser.conf
/etc/daemons.conf
/etc/fstab
/etc/group
/etc/gshadow
/etc/host.conf
/etc/hostname
/etc/hosts
/etc/init.d/rc.functions
/etc/issue
/etc/motd
/etc/networks
/etc/nsswitch.conf
/etc/passwd
/etc/profile
/etc/securetty
/etc/services
/etc/shadow
/etc/shells
/etc/slitaz-release
/dev/core
/dev/flash
/dev/ram
/dev/vcs
/dev/vcsa
/dev/fd
/dev/stdin
/dev/stdout
/dev/stderr
/etc/mtab
