# OSimages/SlitazOS/etc/udev/rules.d/92-hal.rules
# pass all events to hal daemon
RUN+="socket:@/org/freedesktop/hal/udev_event"
