# OSimages/SlitazOS/var/lib/tazpkg/installed/slitaz-tools/files.list
/etc/X11/screen.conf
/etc/firewall.conf
/etc/init.d/crond
/etc/init.d/firewall
/etc/slitaz/applications.conf
/etc/locale.conf
/etc/keymap.conf
/etc/TZ
/sbin/tazlocale
/sbin/tazkeymap
/sbin/tazhw
/sbin/hwsetup
/usr/bin/tazx
/usr/bin/startx
/usr/bin/history
/usr/bin/editor
/usr/bin/browser
/usr/bin/terminal
/usr/bin/file-manager
/usr/sbin/soundconf
/usr/sbin/setmixer
/usr/sbin/tazinst
/usr/share/licenses/COPYING
/usr/share/licenses/bsd.txt
/usr/share/licenses/gpl.txt
/usr/share/licenses/lgpl.txt
/usr/share/licenses/mit.txt
/usr/share/licenses/mozilla.txt
/usr/share/locale/es_AR/LC_MESSAGES/tazinst.mo
/usr/share/locale/es_AR/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/fr/LC_MESSAGES/tazinst.mo
/usr/share/locale/fr/LC_MESSAGES/slitaz-tools.mo
/usr/share/locale/pt_BR/LC_MESSAGES/tazinst.mo
/usr/share/locale/pt_BR/LC_MESSAGES/slitaz-tools.mo
