# OSimages/SlitazOS/var/lib/tazpkg/installed/busybox/files.list
/var/www/cgi-env.sh
/var/www/index.html
/var/www/cgi-bin/index.cgi
/var/www/style.css
/etc/httpd.conf
/etc/dnsd.conf
/etc/zcip.script
/etc/resolv.conf
/etc/udhcpd.conf
/etc/init.d/daemon
/etc/busybox.conf
/etc/inetd.conf
/usr/bin/httpd_helper.sh
/usr/share/applications/httpd-cgi.desktop
/usr/share/applications/httpd.desktop
/usr/share/applications/vi.desktop
/usr/share/udhcpc/default.script
/bin/busybox
/etc/init.d/zcip
/etc/init.d/tftpd
/etc/init.d/syslogd
/etc/init.d/klogd
/etc/init.d/crond
/etc/init.d/dnsd
/etc/init.d/ntpd
/etc/init.d/inetd
/etc/init.d/httpd
/etc/init.d/ftpd
/etc/init.d/udhcpd
/sbin/mkfs.vfat
/sbin/depmod
/sbin/zcip
/sbin/route
/sbin/sysctl
/sbin/tune2fs
/sbin/loadkmap
/sbin/adjtimex
/sbin/getty
/sbin/losetup
/sbin/hdparm
/sbin/hwclock
/sbin/syslogd
/sbin/klogd
/sbin/pivot_root
/sbin/fdisk
/sbin/raidautorun
/sbin/mkdosfs
/sbin/halt
/sbin/switch_root
/sbin/init
/sbin/swapoff
/sbin/ifenslave
/sbin/sulogin
/sbin/blockdev
/sbin/nameif
/sbin/tunctl
/sbin/ifconfig
/sbin/insmod
/sbin/poweroff
/sbin/watchdog
/sbin/reboot
/sbin/freeramdisk
/sbin/blkid
/sbin/findfs
/sbin/arp
/sbin/rmmod
/sbin/mkswap
/sbin/lsmod
/sbin/vconfig
/sbin/udhcpc
/sbin/modinfo
/sbin/swapon
/sbin/fbsplash
/sbin/modprobe
/sbin/setconsole
/usr/sbin/readprofile
/usr/sbin/dhcprelay
/usr/sbin/fakeidentd
/usr/sbin/nbd-client
/usr/sbin/ubidetach
/usr/sbin/crond
/usr/sbin/dnsd
/usr/sbin/setlogcons
/usr/sbin/brctl
/usr/sbin/ubiattach
/usr/sbin/ntpd
/usr/sbin/inetd
/usr/sbin/sendmail
/usr/sbin/rdate
/usr/sbin/httpd
/usr/sbin/ftpd
/usr/sbin/rdev
/usr/sbin/udhcpd
/usr/sbin/chroot
/usr/sbin/chpasswd
/usr/bin/logger
/usr/bin/lzcat
/usr/bin/pgrep
/usr/bin/tac
/usr/bin/split
/usr/bin/expand
/usr/bin/ipcs
/usr/bin/reset
/usr/bin/chat
/usr/bin/[
/usr/bin/crontab
/usr/bin/uniq
/usr/bin/[[
/usr/bin/sum
/usr/bin/tftpd
/usr/bin/xargs
/usr/bin/nslookup
/usr/bin/who
/usr/bin/yes
/usr/bin/pmap
/usr/bin/cmp
/usr/bin/mesg
/usr/bin/wall
/usr/bin/unzip
/usr/bin/telnet
/usr/bin/uptime
/usr/bin/fuser
/usr/bin/last
/usr/bin/tftp
/usr/bin/patch
/usr/bin/timeout
/usr/bin/script
/usr/bin/cksum
/usr/bin/pkill
/usr/bin/mkfifo
/usr/bin/bunzip2
/usr/bin/clear
/usr/bin/ether-wake
/usr/bin/time
/usr/bin/chrt
/usr/bin/killall
/usr/bin/kbd_mode
/usr/bin/renice
/usr/bin/comm
/usr/bin/smemcap
/usr/bin/setsid
/usr/bin/uuencode
/usr/bin/head
/usr/bin/nohup
/usr/bin/fdformat
/usr/bin/cryptpw
/usr/bin/lspci
/usr/bin/sha512sum
/usr/bin/beep
/usr/bin/mkpasswd
/usr/bin/arping
/usr/bin/lsusb
/usr/bin/sha1sum
/usr/bin/dos2unix
/usr/bin/tr
/usr/bin/cut
/usr/bin/cal
/usr/bin/diff
/usr/bin/which
/usr/bin/du
/usr/bin/strings
/usr/bin/traceroute
/usr/bin/vlock
/usr/bin/install
/usr/bin/uudecode
/usr/bin/length
/usr/bin/dpkg-deb
/usr/bin/free
/usr/bin/pscan
/usr/bin/microcom
/usr/bin/tee
/usr/bin/chvt
/usr/bin/top
/usr/bin/wget
/usr/bin/xzcat
/usr/bin/hexdump
/usr/bin/basename
/usr/bin/ftpget
/usr/bin/seq
/usr/bin/whoami
/usr/bin/ipcrm
/usr/bin/realpath
/usr/bin/id
/usr/bin/dpkg
/usr/bin/dumpleases
/usr/bin/logname
/usr/bin/find
/usr/bin/wc
/usr/bin/unxz
/usr/bin/tail
/usr/bin/env
/usr/bin/md5sum
/usr/bin/less
/usr/bin/expr
/usr/bin/unlzma
/usr/bin/hostid
/usr/bin/openvt
/usr/bin/lpr
/usr/bin/passwd
/usr/bin/rtcwake
/usr/bin/sort
/usr/bin/fold
/usr/bin/sha256sum
/usr/bin/nmeter
/usr/bin/taskset
/usr/bin/awk
/usr/bin/bzcat
/usr/bin/deallocvt
/usr/bin/test
/usr/bin/od
/usr/bin/setkeycodes
/usr/bin/dirname
/usr/bin/unix2dos
/usr/bin/flock
/usr/bin/unexpand
/usr/bin/ftpput
/usr/bin/nc
/usr/bin/printf
/usr/bin/tty
/usr/bin/eject
/usr/bin/readlink
/usr/bin/rpm2cpio
/usr/bin/killall5
/bin/iplink
/bin/dd
/bin/mount
/bin/nice
/bin/cp
/bin/gunzip
/bin/sh
/bin/cpio
/bin/mountpoint
/bin/ps
/bin/df
/bin/iostat
/bin/stty
/bin/mv
/bin/iptunnel
/bin/powertop
/bin/kill
/bin/deluser
/bin/login
/bin/chattr
/bin/lsattr
/bin/true
/bin/egrep
/bin/pipe_progress
/bin/iproute
/bin/touch
/bin/date
/bin/mknod
/bin/setarch
/bin/su
/bin/mkdir
/bin/ls
/bin/addgroup
/bin/mktemp
/bin/linux64
/bin/dmesg
/bin/chgrp
/bin/zcat
/bin/rmdir
/bin/mpstat
/bin/cat
/bin/pidof
/bin/adduser
/bin/rev
/bin/conspy
/bin/ipcalc
/bin/sync
/bin/rpm
/bin/ping6
/bin/iprule
/bin/printenv
/bin/fsync
/bin/tar
/bin/chown
/bin/vi
/bin/watch
/bin/rm
/bin/scriptreplay
/bin/gzip
/bin/delgroup
/bin/chmod
/bin/ip
/bin/pwd
/bin/ping
/bin/ionice
/bin/ipaddr
/bin/grep
/bin/ash
/bin/stat
/bin/ln
/bin/uname
/bin/hostname
/bin/echo
/bin/sed
/bin/sleep
/bin/fgrep
/bin/more
/bin/netstat
/bin/dumpkmap
/bin/fdflush
/bin/false
/bin/umount
/bin/linux32
/bin/dnsdomainname
/bin/usleep
/bin/uncompress
